# dr/roles/dr_warm_standby/templates/dr_api.service.j2
[Unit]
Description=DR FastAPI Service (Warm Standby)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ dr_fastapi_user }}
WorkingDirectory={{ dr_fastapi_app_dir }}
ExecStart=/usr/bin/python3 -m uvicorn app:app --host {{ dr_fastapi_host }} --port {{ dr_fastapi_port }}
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
